# neron.service
[Unit]
Description=Néron AI Assistant v2
After=network.target ollama.service
Wants=network-online.target

[Service]
Type=simple
User=__NERON_USER__
Group=__NERON_USER__
WorkingDirectory=__NERON_DIR__/modules/neron_core

# Charger le .env
EnvironmentFile=__NERON_DIR__/.env

# Lancer via le venv
ExecStart=__NERON_DIR__/venv/bin/python3 app.py

# Restart automatique
Restart=always
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=3

# Logs
StandardOutput=append:__NERON_DIR__/logs/neron_core.log
StandardError=append:__NERON_DIR__/logs/neron_core.log

# Timeout
TimeoutStartSec=120
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
