# neron.service
[Unit]
Description=Néron AI Assistant v2
After=network.target ollama.service
Wants=network-online.target

[Service]
Type=simple
User=eleazar
Group=eleazar
WorkingDirectory=/mnt/usb-storage/neron/server

# Charger le .env si besoin
EnvironmentFile=/mnt/usb-storage/neron/server/.env

# Lancer via le venv
ExecStart=/bin/bash -c "source /mnt/usb-storage/neron/server/venv/bin/activate && exec python core/app.py"

# Restart automatique
Restart=always
RestartSec=10

# Timeout
TimeoutStartSec=120
TimeoutStopSec=30

# Logs systemd (optionnel)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
