# neron.service
[Unit]
Description=Néron AI Assistant v2
After=network.target ollama.service
Wants=network-online.target

[Service]
Type=simple
User=__NERON_USER__
Group=__NERON_USER__
WorkingDirectory=__NERON_DIR__

ExecStart=/bin/bash -c "source __NERON_DIR__/venv/bin/activate && exec python core/app.py"

Restart=always
RestartSec=10

TimeoutStartSec=120
TimeoutStopSec=30

StandardOutput=append:__NERON_DIR__/../data/logs/neron_core.log
StandardError=append:__NERON_DIR__/../data/logs/neron_core.log

Environment=PYTHONPATH=__NERON_DIR__/..

[Install]
WantedBy=multi-user.target
