# deploy/systemd/neron-llm.service
[Unit]
Description=Néron LLM microservice
After=network.target

[Service]
Type=simple
User=neron
Group=neron-dev
WorkingDirectory=/etc/neron
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/etc/neron/llm
ExecStart=/etc/neron/venv/bin/uvicorn llm.main:app --host 0.0.0.0 --port 8765 --workers 1
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
