# deploy/neron-llm.service
[Unit]
Description=Néron LLM microservice
After=network.target

[Service]
Type=simple
User=eleazar
WorkingDirectory=/etc/neron/llm
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/etc/neron/llm
ExecStart=/etc/neron/venv/bin/python -m uvicorn neron_llm.main:app --host 0.0.0.0 --port 8765 --workers 1
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
